# build/vmi/modm/files/modm.service
[Unit]
Description=MODM ServiceHost
After=docker.service
StartLimitInterval=200
StartLimitBurst=5

[Service]
Type=notify
EnvironmentFile=/etc/environment
WorkingDirectory=/usr/local/modm/service
ExecStart=/usr/sbin/modm
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target